# site/site_tunnels/templates/etc/systemd/system/tunnel.service.erb
[Unit]
Description=<%= @data['name'] %> tunnel
Requires=network.target
Wants=nss-lookup.target
Before=nss-lookup.target
After=network.target

[Service]
User=tunnels
ExecStart=/usr/bin/ssh -i /etc/tunnels/.ssh/id_<%= @data['name'] %> -N -o "ExitOnForwardFailure=yes" -o "StrictHostKeyChecking=accept-new" -o "ServerAliveInterval=55" <%= @data['forwards'] %> <%= @data['destination'] %> -Nn
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
