# site/profile/templates/etc/systemd/system/wireguard-tunnels.service.epp
[Unit]
Description=Start wireguard tunnels after reboot
Before=bird.service
After=network-online.target nss-lookup.target
Wants=network-online.target nss-lookup.target

[Service]
Type=oneshot
RemainAfterExit=yes
<% $profile::networking::dn42::wireguard::peerings.each |Integer $index, $peer| { -%>
ExecStart=/sbin/ifup wg<%= $index %>
<% } -%>

[Install]
WantedBy=multi-user.target
